# tailscale/rootfs/etc/s6-overlay/s6-rc.d/post-tailscaled/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs after the machine has been logged in into the Tailscale network
# ==============================================================================
declare -a addresses=()
declare -a routes=()
declare ipinfo
declare tags

function appendarray() {
  local -n array=${1}
  readarray -t -O "${#array[@]}" array
}

# Find interfaces and matching addresses from which we can extract routes to be advertised
for interface in $(bashio::network.interfaces); do
  appendarray addresses < <(bashio::network.ipv4_address "${interface}")
  appendarray addresses < <(bashio::network.ipv6_address "${interface}")
done

# Extract routes to be advertised
for address in "${addresses[@]}"; do
  if bashio::var.has_value "${address}"; then
    # Skip local link addresses
    if [[ "${address:0:6}" == "fe80::" ]] || [[ "${address:0:8}" == "169.254." ]];
    then
      continue
    fi

    # Skip if forwarding for the address family is disabled
    if [[ "${address}" =~ .*:.* ]];
    then
      [[ $(</proc/sys/net/ipv6/conf/all/forwarding) -eq 0 ]] && continue
    else
      [[ $(</proc/sys/net/ipv4/ip_forward) -eq 0 ]] && continue
    fi

    # Extract unique routes
    ipinfo="$(/usr/bin/ipcalc --json "${address}")"
    routes+=("$(bashio::jq "${ipinfo}" '.NETWORK + "/" + .PREFIX')")
  fi
done
# Remove duplicate entries
readarray -t routes < <(printf "%s\n" "${routes[@]}" | sort -u)

# Get configured tags
tags=$(bashio::config "tags//[] | join(\",\")" "")

# Wait for socket to be available
while ! bashio::fs.socket_exists "/var/run/tailscale/tailscaled.sock";
do
  sleep 1
done

# Wait for the network to be available
while ! /opt/tailscale status --json --peers=false --self=false \
  | jq --exit-status '.BackendState == "Running" or .BackendState == "NeedsLogin"' > /dev/null;
do
  sleep 2
done

# Start Tailscale
IFS=","
if ! /opt/tailscale up \
  --hostname "$(bashio::info.hostname)" \
  --advertise-exit-node \
  --accept-routes \
  --advertise-routes="${routes[*]}" \
  --ssh \
  --advertise-tags="${tags}";
then
  bashio::log.error "Unable to start up Tailscale"
  bashio::exit.nok
fi

# Wait for the network to be available and logged in
while ! /opt/tailscale status --json --peers=false --self=false \
  | jq --exit-status '.BackendState == "Running"' > /dev/null
do
  sleep 2
done

bashio::log.info "Tailscale is running"
bashio::log.info "Tailscale version: $(/opt/tailscale version --daemon=false | head -n 1)"

bashio::log.info "Enabling Tailscale Proxy..."

declare domain

# Remove previously configured proxy first, because the configuration is persistent between reboots
/opt/tailscale serve https:443 / http://localhost:$(bashio::core.port) off &> /dev/null || true

# Check if Tailscale HTTPS is enabled
if ! /opt/tailscale status --self=true --peers=false --json \
  | jq -rce ".CertDomains[0]" > /dev/null;
then
  bashio::log.info "Can't get your Tailscale domain, is Tailscale's HTTPS support enabled?"
  bashio::log.info "If you want to use https certificates and https proxy on this device,"
  bashio::log.info "please change configuration and restart this add-on."
  bashio::log.info "See the 'Prerequisites' chapter in the add-on's Documentation."
  bashio::exit.ok
fi
domain=$(/opt/tailscale status --self=true --peers=false --json | jq -rc ".CertDomains[0]")

# Checking if SSL is used
if bashio::var.true "$(bashio::core.ssl)"; then
  bashio::log.warning "Tailscale's HTTPS support is enabled, but Home Assistant"
  bashio::log.warning "is not accessible through plain http connection."
  bashio::log.warning "If you wan't to use https certificates and https proxy on this device,"
  bashio::log.warning "please change configuration and restart this add-on."
  bashio::log.warning "See the 'Prerequisites' chapter in the add-on's Documentation."
  bashio::exit.nok
fi

# Set up proxy
if ! /opt/tailscale serve https:443 / http://localhost:$(bashio::core.port); then
  bashio::log.error "Unable to configure Tailscale Proxy"
  bashio::exit.nok
fi
bashio::log.info "Tailscale Proxy is enabled:"
bashio::log.info "  Your Home Assistant instance is available within your Tailnet VPN at"
bashio::log.info "  https://${domain}"

bashio::log.info "Enabling Tailscale Funnel..."

# Check if funnel is available
if ! /opt/tailscale funnel 443 off; then
  bashio::log.info "Funnel is not available, is Tailscale's Funnel feature enabled?"
  bashio::log.info "If you wan't to access this device through Tailscale's Funnel,"
  bashio::log.info "please change configuration and restart this add-on."
  bashio::log.info "See the 'Required Tailscale configuration' chapter in the add-on's Documentation."
  bashio::exit.ok
fi

# Set up funnel
if ! /opt/tailscale funnel 443 on; then
  bashio::log.error "Unable to configure Tailscale Funnel"
  bashio::exit.nok
fi
bashio::log.info "Tailscale Funnel is enabled:"
bashio::log.info "  Your Homeassistant instance is publicly available on the internet at"
bashio::log.info "  https://${domain}"
